# systemd/abridge-watcher.path
[Path]
PathModified=/home/droundy/abridge

[Install]
WantedBy=multi-user.target
